# service/adm-tv/run
#!/bin/bash

cd /www/adm/a6tv
exec dashing start
